# dot_config/systemd/user/lmstudio.service
[Unit]
Description=LM Studio Server
After=network.target

[Service]
Type=forking
ExecStart=/bin/bash -lc 'lms daemon up && lms server start --bind 0.0.0.0 --port 1234'
ExecStop=/bin/bash -lc 'lms server stop && lms daemon down'
Restart=on-failure
RestartSec=10
RemainAfterExit=yes

[Install]
WantedBy=default.target
